# luci-theme-butongwifi/root/etc/uci-defaults/luci-theme-butongwifi
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Butongwifi=/luci-static/butongwifi
	set luci.main.mediaurlbase=/luci-static/butongwifi
	commit luci
EOF
exit 0
